# docker/knot-web/knot-web.service
[Unit]
Description=KNoT Web
Requires=knotd.service

[Service]
ExecStart=/usr/local/bin/knot-web.sh

[Install]
WantedBy=multi-user.target