# docker/knot-web/knot-web.service
[Unit]
Description=KNoT Web

[Service]
ExecStart=/usr/local/bin/knot-web.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target