# docker/wpantund/wpantund.service
[Unit]
Description=WPAN Tunnel Driver

[Service]
ExecStart=/usr/local/bin/wpantund.sh

[Install]
WantedBy=multi-user.target