# docker/knotd/knotd.service
[Unit]
Description=KNoT Daemon

[Service]
ExecStart=/usr/local/bin/knotd.sh

[Install]
WantedBy=multi-user.target